# docker/rhq/root/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700                    

[supervisord]
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run//supervisor.sock

[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/%(program_name)s.log
stderr_logfile=/var/log/%(program_name)s.log
autorestart=true

[program:csassandra]
command=java -ea -javaagent:./../lib/jamm-0.2.5.jar -XX:+UseThreadPriorities -XX:ThreadPriorityPolicy=42 -Xms512M -Xmx512M -Xmn128M -XX:+HeapDumpOnOutOfMemoryError -Xss240k -XX:+UseParNewGC -XX:+UseConcMarkSweepGC -XX:+CMSParallelRemarkEnabled -XX:SurvivorRatio=8 -XX:MaxTenuringThreshold=1 -XX:CMSInitiatingOccupancyFraction=75 -XX:+UseCMSInitiatingOccupancyOnly -XX:+UseTLAB -Djava.rmi.server.hostname=localhost -Dcom.sun.management.jmxremote.port=7299 -Dcom.sun.management.jmxremote.ssl=false -Dcom.sun.management.jmxremote.authenticate=false -Dlog4j.configuration=log4j-server.properties -Dlog4j.defaultInitOverride=true -Dcassandra-pidfile=/opt/rhq/rhq-storage/bin/cassandra.pid -cp ./../conf:./../build/classes/main:./../build/classes/thrift:./../lib/antlr-3.2.jar:./../lib/apache-cassandra-1.2.9.jar:./../lib/apache-cassandra-clientutil-1.2.9.jar:./../lib/apache-cassandra-thrift-1.2.9.jar:./../lib/avro-1.4.0-fixes.jar:./../lib/avro-1.4.0-sources-fixes.jar:./../lib/commons-cli-1.1.jar:./../lib/commons-codec-1.2.jar:./../lib/commons-lang-2.6.jar:./../lib/compress-lzf-0.8.4.jar:./../lib/concurrentlinkedhashmap-lru-1.3.jar:./../lib/guava-13.0.1.jar:./../lib/high-scale-lib-1.1.2.jar:./../lib/jackson-core-asl-1.9.2.jar:./../lib/jackson-mapper-asl-1.9.2.jar:./../lib/jamm-0.2.5.jar:./../lib/jbcrypt-0.3m.jar:./../lib/jline-1.0.jar:./../lib/json-simple-1.1.jar:./../lib/libthrift-0.7.0.jar:./../lib/log4j-1.2.16.jar:./../lib/metrics-core-2.0.3.jar:./../lib/netty-3.6.6.Final.jar:./../lib/rhq-cassandra-auth-4.9.0.jar:./../lib/servlet-api-2.5-20081211.jar:./../lib/slf4j-api-1.7.2.jar:./../lib/slf4j-log4j12-1.7.2.jar:./../lib/snakeyaml-1.6.jar:./../lib/snappy-java-1.0.4.1-rhq-p1.jar:./../lib/snaptree-0.1.jar org.apache.cassandra.service.CassandraDaemon
directory=/opt/rhq/rhq-storage/bin
stdout_logfile=/var/log/%(program_name)s.log
stderr_logfile=/var/log/%(program_name)s.log

[program:rhq-server]
command=/usr/lib/jvm/java-7-openjdk-amd64/jre/bin/java -D[Standalone] -XX:+UseCompressedOops -Dapp.name=rhq-server -Xms1024M -Xmx1024M -XX:PermSize=256M -XX:MaxPermSize=256M -Djava.net.preferIPv4Stack=true -Dorg.jboss.resolver.warning=true -Dsun.rmi.dgc.client.gcInterval=3600000 -Dsun.rmi.dgc.server.gcInterval=3600000 -Drhq.server.home=/opt/rhq -Djboss.server.log.dir=/opt/rhq/logs -Djava.awt.headless=true -Dsun.lang.ClassLoader.allowArraySyntax=true -Djboss.server.default.config=standalone-full.xml -Djboss.modules.system.pkgs=org.jboss.byteman -Djava.security.manager -Djava.security.policy==/opt/rhq/bin/rhq-server.security-policy -Dorg.apache.coyote.http11.Http11Protocol.COMPRESSION=on -Dorg.apache.coyote.http11.Http11Protocol.COMPRESSION_MIME_TYPES=text/javascript,text/css,text/html -Dorg.jboss.boot.log.file=/opt/rhq/logs/server.log -Dlogging.configuration=file:/opt/rhq/jbossas/standalone/configuration/logging.properties -jar /opt/rhq/jbossas/jboss-modules.jar -mp /opt/rhq/modules:/opt/rhq/jbossas/modules -jaxpmodule javax.xml.jaxp-provider org.jboss.as.standalone -Djboss.home.dir=/opt/rhq/jbossas -Djboss.server.base.dir=/opt/rhq/jbossas/standalone -P /opt/rhq/bin/rhq-server.properties
stdout_logfile=/var/log/%(program_name)s.log
stderr_logfile=/var/log/%(program_name)s.log

[program:rhq-agent]
command=/usr/bin/java -Djava.endorsed.dirs=/opt/rhq-agent/lib/endorsed -Djava.library.path=/opt/rhq-agent/lib -Djna.platform.library.path=/opt/rhq-agent/lib/augeas/lib64 -Xms64m -Xmx128m -Djava.net.preferIPv4Stack=true -Drhq.preferences.file=/opt/rhq-agent/conf/agent-prefs.properties -Dlog4j.configuration=log4j.xml -cp /opt/rhq-agent/conf:/opt/rhq-agent/lib/commons-io-1.4.jar:/opt/rhq-agent/lib/commons-logging-1.1.0.jboss.jar:/opt/rhq-agent/lib/concurrent-1.3.4-jboss-update1.jar:/opt/rhq-agent/lib/getopt-1.0.13.jar:/opt/rhq-agent/lib/i18nlog-1.0.10.jar:/opt/rhq-agent/lib/jboss-common-core-2.2.17.GA.jar:/opt/rhq-agent/lib/jboss-jmx-4.2.3.GA.jar:/opt/rhq-agent/lib/jboss-logging-3.1.2.GA.jar:/opt/rhq-agent/lib/jboss-remoting-2.5.4.SP4.jar:/opt/rhq-agent/lib/jboss-serialization-1.0.3.GA.jar:/opt/rhq-agent/lib/jline-0.9.94.jar:/opt/rhq-agent/lib/log4j-1.2.16.jar:/opt/rhq-agent/lib/persistence-api-1.0.jar:/opt/rhq-agent/lib/rhq-common-drift-4.9.0.jar:/opt/rhq-agent/lib/rhq-core-client-api-4.9.0.jar:/opt/rhq-agent/lib/rhq-core-comm-api-4.9.0.jar:/opt/rhq-agent/lib/rhq-core-domain-4.9.0.jar:/opt/rhq-agent/lib/rhq-core-native-system-4.9.0.jar:/opt/rhq-agent/lib/rhq-core-plugin-api-4.9.0.jar:/opt/rhq-agent/lib/rhq-core-plugin-container-4.9.0.jar:/opt/rhq-agent/lib/rhq-core-util-4.9.0.jar:/opt/rhq-agent/lib/rhq-enterprise-agent-4.9.0.jar:/opt/rhq-agent/lib/rhq-enterprise-comm-4.9.0.jar:/opt/rhq-agent/lib/sigar-1.6.5.132-5.jar org.rhq.enterprise.agent.AgentMain
stdout_logfile=/var/log/%(program_name)s.log
stderr_logfile=/var/log/%(program_name)s.log
