# docker/as7/root/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700   

[inet_http_server] 
port = *:9001
;username=webstack
;password=secret         

[supervisord]
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run//supervisor.sock

[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/%(program_name)s.log
stderr_logfile=/var/log/%(program_name)s.log
autorestart=true

[program:wildfly]
command=%(ENV_JAVA)s
    -D[Standalone] 
    %(ENV_JAVA_OPTS)s
    -Djava.net.preferIPv4Stack=true
    -Djboss.modules.system.pkgs=org.jboss.byteman
    -Djava.awt.headless=true
    -Dorg.jboss.boot.log.file=%(ENV_WILDFLY_LOG_DIR)s/wildfly.boot.log
    -Dlogging.configuration=file:%(ENV_WILDFLY_CONFIG_DIR)s/logging.properties
    -jar "%(ENV_WILDFLY_HOME)s/jboss-modules.jar"
    -mp "%(ENV_WILDFLY_MODULE_DIR)s"
    -jaxpmodule "javax.xml.jaxp-provider"
    org.jboss.as.standalone
    -Djboss.home.dir="%(ENV_WILDFLY_HOME)s"
    -Djboss.server.base.dir="%(ENV_WILDFLY_BASE_DIR)s"

stdout_logfile=/var/log/%(program_name)s.log
stderr_logfile=/var/log/%(program_name)s.error.log
user=wildfly
